# deploy/alumni-book-cleanup.service
[Unit]
Description=Alumni Book expired session cleanup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/alumni-book/app
ExecStart=/usr/bin/podman-compose -f /opt/alumni-book/app/docker-compose.yml exec -T api pnpm --filter worker cleanup:local
