# deploy/alumni-book-backup.service
[Unit]
Description=Alumni Book SQLite and upload backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
Group=admin
Environment=XDG_RUNTIME_DIR=/run/user/1001
WorkingDirectory=/opt/alumni-book/app
ExecStart=/usr/bin/bash /opt/alumni-book/app/scripts/backup-selfhosted.sh
