# deploy/alumni-book-backup.service
[Unit]
Description=Alumni Book SQLite and upload backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/alumni-book/app
ExecStart=/opt/alumni-book/app/scripts/backup-selfhosted.sh
