# deploy/alumni-book-cleanup.service
[Unit]
Description=Alumni Book expired session cleanup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
Group=admin
Environment=XDG_RUNTIME_DIR=/run/user/1001
WorkingDirectory=/opt/alumni-book/app
ExecStart=/usr/bin/podman-compose -f /opt/alumni-book/app/docker-compose.yml exec -T api pnpm --filter worker cleanup:local
