# deploy/alumni-book-backup.timer
[Unit]
Description=Run Alumni Book backup daily

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
